# deploy/cron.daily
# Run 00:15 UTC — one line per minute, overwrites /etc/cron.d/forgeone-daily.
# Install with: sudo install -m 0644 deploy/cron.daily /etc/cron.d/forgeone-daily
15 0 * * * root cd /opt/forgeone && /opt/forgeone/.venv/bin/python -m forgeone.scripts.replay_paper_day --post-relay >> /opt/forgeone/data/daily_reports/_cron.log 2>&1
